# opentsdb/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:hbase]
command=/hbase-1.2.1/bin/start-hbase.sh
stopwaitsecs=60

[program:opentsdb]
command=/start.sh

